# overlay/system/etc/systemd/system/sensorfwd.service.d/wait-sensors-hal.conf
[Service]
ExecStartPre=/bin/bash -c 'while [ "$(getprop init.svc.vendor.sensors-hal-2-0-multihal)" != "running" ]; do sleep 0.5; done'
